# knowledge/infra/openshift/config/rc.local
touch /var/lock/subsys/local

/home/azadmin/cmd/mount-aznfs.sh